# scripts/ops/supervisor-queue.conf
[program:amigo-secreto-queue]
process_name=%(program_name)s_%(process_num)02d

; Prioridade: processa a fila 'emails' antes da 'default'.
; Isso garante que e-mails de verificação de cadastro não fiquem atrás
; de outros jobs (ex.: sorteios) na fila de espera.
;
; --timeout=60: Menor que retry_after=90 (config/queue.php) para evitar
;               que dois workers processem o mesmo job simultaneamente.
;
; --max-jobs=500: Reinicia o worker preventivamente a cada 500 jobs
;                 para mitigar memory leaks em execuções longas.
;
; --backoff=5: Aguarda 5s antes de tentar reprocessar um job com falha.
command=php /var/www/amigo-secreto/artisan queue:work --queue=emails,default --sleep=3 --tries=3 --timeout=60 --backoff=5 --max-jobs=500

autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/amigo-secreto/storage/logs/queue-worker.log
stopwaitsecs=3600
